# crates/ruptela-listener/debian/ruptela-listener.service
[Unit]
Description=Ruptela GPS TCP Listener
After=network.target

[Service]
Type=simple
User=ruptela
LimitNOFILE=65536
EnvironmentFile=/etc/ruptela/listener
ExecStart=/usr/ruptela/listener \
    --host ${LISTENER_HOST} \
    --port ${LISTENER_PORT} \
    --redis-url ${REDIS_URL} \
    --shards ${SHARDS} \
    --max-connections ${MAX_CONNECTIONS}
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
